# modules/quadlets/post-boot/user-quadlets-setup.timer
[Unit]
Description=Timer for user quadlets setup

[Timer]
OnBootSec=5m
OnUnitInactiveSec=24h

[Install]
WantedBy=timers.target
