# files/usr/lib/systemd/user/flatpaksync.path
[Unit]
Description=Monitorar o diretório ~/.var/app

[Path]
PathChanged=%h/.var/app
PathChanged=%h/.config/flatpaks.user.installed
Unit=flatpaksync.service

[Install]
WantedBy=default.target
