# files/usr/lib/systemd/user/flatpaksync.path
[Unit]
Description=Service monitoring the ~/.var/app directory

[Path]
PathChanged=%h/.var/app
PathChanged=%h/.config/flatpaks.user.installed
Unit=flatpaksync.service

[Install]
WantedBy=default.target
