# modules/default-flatpaks/v2/post-boot/system-flatpak-setup.timer
[Unit]
Description=Timer for system-flatpak-setup

[Timer]
OnBootSec=30

[Install]
WantedBy=timers.target