# modules/default-flatpaks/v2/post-boot/user-flatpak-setup.timer
[Unit]
Description=Timer for user-flatpak-setup

[Timer]
OnBootSec=30

[Install]
WantedBy=timers.target